# deploy/supervisor-app.conf
[supervisord]
nodaemon=true

[program:app-srs]
command = /opt/srs/objs/srs -c /opt/srs/conf/srs.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nginx-app]
command = /usr/sbin/nginx
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
